# rootfs/etc/hotplug.d/scsi_disk/00-udisk-check.sh
#!/bin/sh

[ "${ACTION}" = "add" ] && touch /var/run/udisk-upgrade-pending

